# vault-dashboard/scripts/com.revivr.asc-sync-ids.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Reference copy. Install with:
       cp scripts/com.revivr.asc-sync-ids.plist ~/Library/LaunchAgents/
       launchctl load -w ~/Library/LaunchAgents/com.revivr.asc-sync-ids.plist
     Weekly (Mon 06:30) reconcile of App Store IDs into vault profiles. App IDs
     change rarely, so weekly is plenty. Points at the RUNTIME script. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.revivr.asc-sync-ids</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/einarjohnson/Library/Application Support/Revivr/VaultDashboard/scripts/asc-sync-ids.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/einarjohnson/Library/Logs/revivr-asc-sync-ids.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/einarjohnson/Library/Logs/revivr-asc-sync-ids.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/einarjohnson</string>
    </dict>
</dict>
</plist>
