# vault-dashboard/scripts/com.revivr.awareness-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Reference copy. Install with:
       cp scripts/com.revivr.awareness-refresh.plist ~/Library/LaunchAgents/
       launchctl load -w ~/Library/LaunchAgents/com.revivr.awareness-refresh.plist
     Fetches feeds + distills the daily briefing at 07:00. Points at the
     RUNTIME script (deployed by deploy-dashboard.sh). -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.revivr.awareness-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/einarjohnson/Library/Application Support/Revivr/VaultDashboard/scripts/awareness-refresh.sh</string>
    </array>

    <!-- Daily at 07:00 — one fetch + briefing per day. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/einarjohnson/Library/Logs/revivr-awareness-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/einarjohnson/Library/Logs/revivr-awareness-refresh.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/einarjohnson</string>
    </dict>
</dict>
</plist>
